# ecsc/ecsc2025-infra/packer/ansible/roles/vulnbox_artifact_services/files/start-vulnservices.service
[Unit]
Description=Start Services
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/bash /root/services/start.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
